# L4B5/537UP/wifi-check.cron
* * * * * /home/serfary/L1GH7H0U53/L4B5/537UP/wifi-check.sh >/dev/null 2>&1
